# daily-news-watchdog.service
[Unit]
Description=Daily News Watchdog
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/daily-news-watchdog
ExecStart=/opt/daily-news-watchdog/venv/bin/python watchdog.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daily-news-watchdog

[Install]
WantedBy=multi-user.target
